# supervisord.conf
; ============================================================
; Supervisord — Dexter + VNStock process manager
; ============================================================
; Startup order: vnstock (priority 10) starts before dexter (priority 20)
; Both services auto-restart on failure.
; Logs → stdout/stderr so Docker / App Runner captures them.
; ============================================================

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

; ── VNStock FastAPI service ───────────────────────────────────
[program:vnstock]
command=python3 /app/vnstock/vnstock_service.py
directory=/app/vnstock
autostart=true
autorestart=true
; Give FastAPI time to fully initialize before Dexter starts
startsecs=5
startretries=3
priority=10
; Pipe logs to container stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Friendly process name in logs
process_name=vnstock

; ── Dexter HTTP API server (Bun) ─────────────────────────────
[program:dexter]
command=/root/.bun/bin/bun run /app/dexter/src/server.ts
directory=/app/dexter
autostart=true
autorestart=true
; Wait longer — Dexter loads tools/skills at startup
startsecs=10
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
process_name=dexter

; ── Supervisor web UI (optional, disabled in prod) ───────────
; Uncomment below to enable admin UI at :9001 for debugging:
; [inet_http_server]
; port=0.0.0.0:9001
; username=admin
; password=changeme
